# deploy/systemd/cr-backup-db.service
[Unit]
Description=CR — Denní záloha produkční PostgreSQL databáze na Cloudflare R2
Documentation=https://github.com/Olbrasoft/cr (task #97)
# Wait for docker — skript mluví s kontejnerem `cr-db-1` přes `docker compose exec`.
After=docker.service network-online.target
Wants=network-online.target
StartLimitBurst=2
StartLimitIntervalSec=600

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/cr
# /opt/cr/.env drží R2_BACKUP_* credentials (viz ~/Dokumenty/přístupy/cloudflare/r2-cr-db-backups.md).
# Skript si sám exportuje RCLONE_CONFIG_CR_R2_BACKUP_* z R2_BACKUP_* — systemd
# `Environment="...=${FOO}"` neumí shell expansion, takže to nemůže být tady.
EnvironmentFile=/opt/cr/.env
ExecStart=/bin/bash /opt/cr/scripts/backup-db.sh auto
# Loguj do journald (default), ne do souboru bez rotace. journalctl -u stačí.
StandardOutput=journal
StandardError=journal
# Drop škrábané logy (i když 1 log/s ze skriptu by neměl nikdy být problém).
LogRateLimitIntervalSec=10s
LogRateLimitBurst=200
# Pokud transientně selže (Cloudflare 502, DB v restartu), zkus ještě jednou.
# Pokud padne i podruhé, nech to být — admin banner upozorní ráno.
Restart=on-failure
RestartSec=120

[Install]
WantedBy=multi-user.target
