# deploy/systemd/cr-llm-resolver.service
[Unit]
Description=CR — LLM-based TMDB-ID resolver for prehraj.to unmatched clusters (#668)
Documentation=https://github.com/Olbrasoft/cr/issues/652
# Wait for Docker / Postgres to be up. `After=cr-prehrajto-sync` only
# orders start jobs in the same transaction — it does NOT prevent the
# resolver from running while the sync is still active (e.g., when the
# sitemap import takes longer than 2.5 hours, or when this service is
# triggered manually). The ExecStartPre below adds an explicit guard.
After=docker.service network-online.target cr-prehrajto-sync.service
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/cr
EnvironmentFile=/opt/cr/.env
# Refuse to start if the sync is mid-flight — the resolver depends on
# a coherent `prehrajto_unmatched_clusters` registry, and the sync's
# importer step (mark-dead + UPSERT) writes to that table. Running
# concurrently would let the resolver re-attempt clusters the importer
# is in the middle of clearing.
ExecStartPre=/bin/bash -c 'if systemctl is-active --quiet cr-prehrajto-sync.service; then echo "cr-prehrajto-sync.service still active — refusing to start" >&2; exit 1; fi'
# DATABASE_URL in /opt/cr/.env points at `db:5432` (container hostname);
# host-side scripts need 127.0.0.1 — same translation as the
# cr-auto-import / cr-prehrajto-sync services. `exec` replaces the
# bash process with python3 so systemd tracks the actual worker PID
# (matches cr-auto-import.service; better signal/timeout handling
# than a bash parent + python child tree).
#
# --limit 200: caps each daily run at ~7 minutes wall-clock and ~3% of
#   Gemma's 14,400/day free quota. The 10k+ backlog drains over weeks,
#   not days — that's intentional, so a buggy resolver iteration can't
#   burn the entire backlog before we notice false positives.
# --min-uploads 2: ignore single-upload clusters (long-tail noise).
# --retry-after-days 7: don't reprocess clusters the resolver already
#   touched within the last 7 days.
ExecStart=/bin/bash -c 'DATABASE_URL=$${DATABASE_URL//@db:/@127.0.0.1:} exec /usr/bin/python3 -u /opt/cr/scripts/resolve-unmatched-via-llm.py --limit 200 --min-uploads 2 --retry-after-days 7'
StandardOutput=append:/var/log/cr-llm-resolver.log
StandardError=append:/var/log/cr-llm-resolver.log
# No auto-restart: a Gemma quota / TMDB outage shouldn't cause an
# immediate retry that re-enters quota issues. The next daily timer
# tick will pick up where this left off (per-row commits + skip-window
# make this safe).
TimeoutStartSec=30min

[Install]
WantedBy=multi-user.target
